# examples/6lbr/package/DEBIAN/postrm
#!/bin/sh
set -e

if [ "$1" = "purge" ] ; then
	update-rc.d 6lbr remove >/dev/null
fi
