# examples/6lbr/package/DEBIAN/postinst
#!/bin/sh
set -e

action="$1"

if [ "$action" != configure ]; then
    exit 0
fi

if [ -x "/etc/init.d/6lbr" ]; then
	update-rc.d 6lbr defaults >/dev/null
fi
